# build/usehostovs/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=/var/log/supervisord.log                ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=5                               ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                   ; do not run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/                           ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

;
; Runs the ovs forwarder
;
[program:ovs-forwarder]
command=/bin/ovs-forwarder
autorestart=true
startretries=3
stdout_logfile=/var/log/ovs-forwarder.log
stderr_logfile=/var/log/ovs-forwarder.log
priority=30
startsecs=5

;
; Redirect ovs forwarder log to console
;
[program:forwarder-stderr]
command=tail -f /var/log/ovs-forwarder.log
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
